# templates/service/transmart-solr.service.erb
[Unit]
Description=TranSMART Solr service

[Service]
Type=forking
ExecStart=<%= @solr_script %> start
ExecReload=<%= @solr_script %> restart
ExecStop=<%= @solr_script %> stop
User=<%= @user %>
Restart=always
StandardOutput=journal

[Install]
WantedBy=multi-user.target

