# templates/service/transmart-server.service.erb
[Unit]
Description=TranSMART Service

[Service]
ExecStart=<%= @start_script %>
User=<%= @user %>
Restart=always
StandardOutput=journal

[Install]
WantedBy=multi-user.target

