# templates/service/transmart-rserve.service.erb
[Unit]
Description=RServe

[Service]
ExecStart=<%= @rserve_script %>
User=<%= @user %>
Restart=always
StandardOutput=journal

[Install]
WantedBy=multi-user.target

