# templates/transmart-app.service.erb
[Unit]
Description=TranSMART Service

[Service]
ExecStart=<%= @start_script %>
User=<%= @user %>
Restart=always
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target

